# boobytrap.service
[Unit]
Description=Boobytrap

[Service]
Type=simple
User=root
Group=root
ExecStart=boobytrap

[Install]
WantedBy=multi-user.target
