# boobytrap.service
[Unit]
Description=Boobytrap

[Service]
Type=simple
User=root
Group=root
ExecStart=/var/lib/boobytrap/run.sh

[Install]
WantedBy=multi-user.target
